# config/systemd/user/newsboat-refresh.timer
[Unit]
Description=Refresh RSS Feeds Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
