# systemd/user/newsboat-refresh.service
[Unit]
Description=Refresh RSS Feeds

[Service]
Type=oneshot
ExecStart=/usr/bin/newsboat -x reload
